# packaging/systemd/filewall-ui.service
[Unit]
Description=filewall access-prompt UI helper
After=graphical-session.target
PartOf=graphical-session.target
# Only one UI may hold the daemon's single prompt socket; enabling the iced
# helper stops this one (and vice versa).
Conflicts=filewall-ui-iced.service

[Service]
Type=simple
ExecStart=/usr/bin/filewall-ui
Restart=on-failure
RestartSec=2
# yad needs the session's DISPLAY/WAYLAND_DISPLAY + DBUS_SESSION_BUS_ADDRESS.
# On a systemd-managed graphical session the DE imports these into the user
# manager's environment, so graphical-session.target carries them. Bare WMs that
# don't import session env need, typically from ~/.xinitrc:
#   systemctl --user import-environment DISPLAY WAYLAND_DISPLAY XAUTHORITY DBUS_SESSION_BUS_ADDRESS
# filewall-ui fails closed (DenyOnce) if yad can't render, so a missing display
# degrades safely rather than granting access.

[Install]
WantedBy=graphical-session.target
